# images/svtech_icinga2/configs/profile.d/env.sh
export HISTTIMEFORMAT='%F %T  '
export TERM=xterm
umask 0000
